# ansible/roles/qos_per_user_xray/templates/qos-gc.timer.j2
[Unit]
Description=Run VFF QoS tc GC periodically

[Timer]
# Run once per day
OnCalendar=daily

# Add jitter to avoid synchronized load spikes (systemd will delay start randomly)
RandomizedDelaySec=1h

# Ensure missed runs execute soon after boot (if the machine was down at schedule time)
Persistent=true

# Accuracy is not critical for daily jobs; relax wakeups
AccuracySec=5min

Unit=qos-gc.service

[Install]
WantedBy=timers.target
